# flask-microblog/server_config/microblog_gun.conf
[program:gunicorn]
command=gunicorn microblog:app
directory: /home/ubuntu/flask-microblog
user=nobody
autostart=true
autorestart=true
redirect_stderr=true
